# systemd/.config/systemd/user/nvim-clear-lsp.service
[Unit]
Description=Clear Neovim LSP log

[Service]
Type=oneshot
ExecStart=rm -f /home/ddperez/.local/state/nvim/lsp.log
